# rootfs/etc/services.d/httpd/run
#!@execline@/bin/execlineb -P
unexport HOME
unexport PWD
export PATH /bin
@apacheHttpd@/bin/httpd -D FOREGROUND -d @rootfs@/etc/httpd
